# configs/systemd/pam-oauth2-broker.service
[Unit]
Description=pam-oauth2 Authentication Broker
Documentation=https://github.com/scttfrdmn/pam-oauth2
After=network.target
Wants=network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/pam-oauth2-broker --config /etc/pam-oauth2/broker.yaml
Restart=always
RestartSec=5s

# Run as root so it can write to /var/run/pam-oauth2/ and /var/log/pam-oauth2/
User=root
Group=root

# Runtime directory
RuntimeDirectory=pam-oauth2
RuntimeDirectoryMode=0755

# Log directory
LogsDirectory=pam-oauth2
LogsDirectoryMode=0750

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/var/run/pam-oauth2 /var/log/pam-oauth2 /etc/pam-oauth2
PrivateTmp=true

[Install]
WantedBy=multi-user.target
